# ubuntu18/init-system-helpers/debian/init-system-helpers.install
script/deb-systemd-helper usr/bin
script/deb-systemd-invoke usr/bin
script/invoke-rc.d usr/sbin
script/service usr/sbin
script/update-rc.d usr/sbin
